# etc/systemd/system/wm_timechecker.timer
[Unit]
Description=WM timechecker timer

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s

[Install]
WantedBy=timers.target