# etc/systemd/system/wm_timechecker.service
[Unit]
Description=Service for WM time check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/start.sh

[Install]
WantedBy=multi-user.target
